# example/examplebot.conf
description "Example Bot"
author "Someone"

start on startup
stop on shutdown
respawn

script
sudo -u root php -f /path/to/run.php
end script